# _images/devenv-web-desktop/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:pulseaudio]
command=/usr/bin/pulseaudio --daemonize=no --exit-idle-time=-1 --load="module-null-sink sink_name=output sink_properties=device.description=output" --load="module-null-sink sink_name=input sink_properties=device.description=input"
user=developer
environment=HOME="/home/developer",XDG_RUNTIME_DIR="/tmp/developer-run",DBUS_SESSION_BUS_ADDRESS="%(ENV_DBUS_SESSION_BUS_ADDRESS)s"
autorestart=true
priority=200

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 %(ENV_RESOLUTION)s +extension COMPOSITE +extension DAMAGE +extension RANDR +extension RENDER +extension MIT-SHM +extension XFIXES +extension XTEST -nolisten tcp -ac -noreset -shmem
user=developer
environment=HOME="/home/developer"
autorestart=true
priority=300

[program:xfce4]
command=/usr/bin/xfce4-session
user=developer
environment=DISPLAY=":1",HOME="/home/developer",XDG_RUNTIME_DIR="/tmp/developer-run",DBUS_SESSION_BUS_ADDRESS="%(ENV_DBUS_SESSION_BUS_ADDRESS)s"
autorestart=true
priority=400

[program:selkies]
command=/opt/selkies/bin/selkies --addr=localhost --mode=websockets
user=developer
environment=DISPLAY=":1",HOME="/home/developer",SELKIES_ENCODER="%(ENV_SELKIES_ENCODER)s",XDG_RUNTIME_DIR="/tmp/developer-run",DBUS_SESSION_BUS_ADDRESS="%(ENV_DBUS_SESSION_BUS_ADDRESS)s"
autorestart=true
priority=500

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=developer
environment=HOME="/home/developer"
autorestart=true
priority=600
